# rpi/bin/get-sensor.service
[Unit]
Description=Get-Sensor

[Service]
ExecStart=/usr/bin/bash /home/brettonw/bin/get-sensor.sh
Restart=on-abort
User=brettonw
Type=simple

[Install]
WantedBy=multi-user.target

